# scheduller/supervisord_conf/scraper2.supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[inet_http_server]         ; inet (TCP) server disabled by default
port=*:9001        ; (ip_address:port specifier, *:port for all iface)
username=ionscraper                       ; (default is no username (open server))
password=m@nt@bg@n                  ; (default is no password (open server))

[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisord/               ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:crawl_antara]
command=scrapy crawl antara --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-antara-stderr.log
stdout_logfile = /var/log/supervisord/crawl-antara-stdout.log

[program:crawl_bbcindo]
command=scrapy crawl bbcindo --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-bbcindo-stderr.log
stdout_logfile = /var/log/supervisord/crawl-bbcindo-stdout.log

[program:crawl_beritasatu]
command=scrapy crawl beritasatu --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-beritasatu-stderr.log
stdout_logfile = /var/log/supervisord/crawl-beritasatu-stdout.log

[program:crawl_bijaks]
command=scrapy crawl bijaks --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-bijaks-stderr.log
stdout_logfile = /var/log/supervisord/crawl-bijaks-stdout.log

[program:crawl_bisnis]
command=scrapy crawl bisnis --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-bisnis-stderr.log
stdout_logfile = /var/log/supervisord/crawl-bisnis-stdout.log

[program:crawl_jakartaglobe]
command=scrapy crawl jakartaglobe --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-jakartaglobe-stderr.log
stdout_logfile = /var/log/supervisord/crawl-jakartaglobe-stdout.log

[program:crawl_jawapos]
command=scrapy crawl jawapos --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-jawapos-stderr.log
stdout_logfile = /var/log/supervisord/crawl-jawapos-stdout.log

[program:crawl_kontan]
command=scrapy crawl kontan --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-kontan-stderr.log
stdout_logfile = /var/log/supervisord/crawl-kontan-stdout.log

[program:crawl_swa]
command=scrapy crawl swa --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-swa-stderr.log
stdout_logfile = /var/log/supervisord/crawl-swa-stdout.log

[program:crawl_viva]
command=scrapy crawl viva --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-viva-stderr.log
stdout_logfile = /var/log/supervisord/crawl-viva-stdout.log

[program:crawl_thejakartapost]
command=scrapy crawl thejakartapost --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-thejakartapost-stderr.log
stdout_logfile = /var/log/supervisord/crawl-thejakartapost-stdout.log

[program:crawl_bbc]
command=scrapy crawl bbc --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-bbc-stderr.log
stdout_logfile = /var/log/supervisord/crawl-bbc-stdout.log

[program:crawl_nytimes]
command=scrapy crawl nytimes --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-nytimes-stderr.log
stdout_logfile = /var/log/supervisord/crawl-nytimes-stdout.log

[program:crawl_washington]
command=scrapy crawl washington --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-washington-stderr.log
stdout_logfile = /var/log/supervisord/crawl-washington-stdout.log

[program:crawl_detik]
command=scrapy crawl detik --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-detik-stderr.log
stdout_logfile = /var/log/supervisord/crawl-detik-stdout.log

