# scheduller/supervisord_conf/scraper1.supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock                       ; path to your socket file

[inet_http_server]         ; inet (TCP) server disabled by default
port=*:9001        ; (ip_address:port specifier, *:port for all iface)
username=ionscraper                       ; (default is no username (open server))
password=m@nt@bg@n                  ; (default is no password (open server))

[supervisord]
logfile=/var/log/supervisord/supervisord.log    ; supervisord log file
logfile_maxbytes=50MB                           ; maximum size of logfile before rotation
logfile_backups=10                              ; number of backed up logfiles
loglevel=error                                  ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid                ; pidfile location
nodaemon=false                                  ; run supervisord as a daemon
minfds=1024                                     ; number of startup file descriptors
minprocs=200                                    ; number of process descriptors
user=root                                       ; default user
childlogdir=/var/log/supervisord/               ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:crawl_cnnindonesia]
command=scrapy crawl cnnindonesia --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-cnnindonesia-stderr.log
stdout_logfile = /var/log/supervisord/crawl-cnnindonesia-stdout.log

[program:crawl_inilah]
command=scrapy crawl inilah --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-inilah-stderr.log
stdout_logfile = /var/log/supervisord/crawl-inilah-stdout.log

[program:crawl_kompas]
command=scrapy crawl kompas --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-kompas-stderr.log
stdout_logfile = /var/log/supervisord/crawl-kompas-stdout.log

[program:crawl_mediaindonesia]
command=scrapy crawl mediaindonesia --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-mediaindonesia-stderr.log
stdout_logfile = /var/log/supervisord/crawl-mediaindonesia-stdout.log

[program:crawl_merdeka]
command=scrapy crawl merdeka --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-merdeka-stderr.log
stdout_logfile = /var/log/supervisord/crawl-merdeka-stdout.log

[program:crawl_metrotv]
command=scrapy crawl metrotv --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-metrotv-stderr.log
stdout_logfile = /var/log/supervisord/crawl-metrotv-stdout.log

[program:crawl_okezone]
command=scrapy crawl okezone --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-okezone-stderr.log
stdout_logfile = /var/log/supervisord/crawl-okezone-stdout.log

[program:crawl_pikiranrakyat]
command=scrapy crawl pikiranrakyat --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-pikiranrakyat-stderr.log
stdout_logfile = /var/log/supervisord/crawl-pikiranrakyat-stdout.log

[program:crawl_rmol]
command=scrapy crawl rmol --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-rmol-stderr.log
stdout_logfile = /var/log/supervisord/crawl-rmol-stdout.log

[program:crawl_suara]
command=scrapy crawl suara --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-suara-stderr.log
stdout_logfile = /var/log/supervisord/crawl-suara-stdout.log

[program:crawl_smh]
command=scrapy crawl smh --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-smh-stderr.log
stdout_logfile = /var/log/supervisord/crawl-smh-stdout.log

[program:crawl_nationalgeographic]
command=scrapy crawl nationalgeographic --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-nationalgeographic-stderr.log
stdout_logfile = /var/log/supervisord/crawl-nationalgeographic-stdout.log

[program:crawl_aljazeera]
command=scrapy crawl aljazeera --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-aljazeera-stderr.log
stdout_logfile = /var/log/supervisord/crawl-aljazeera-stdout.log

[program:crawl_reuters]
command=scrapy crawl reuters --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-reuters-stderr.log
stdout_logfile = /var/log/supervisord/crawl-reuters-stdout.log

[program:crawl_huffingtonpost]
command=scrapy crawl huffingtonpost --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-huffingtonpost-stderr.log
stdout_logfile = /var/log/supervisord/crawl-huffingtonpost-stdout.log

[program:crawl_editioncnn]
command=scrapy crawl editioncnn --loglevel=INFO
directory=/home/ionscraper
stderr_logfile = /var/log/supervisord/crawl-editioncnn-stderr.log
stdout_logfile = /var/log/supervisord/crawl-editioncnn-stdout.log

